# ansible/roles/polkadot-public/templates/polkadot.service.j2
[Unit]
Description=Polkadot Node

[Service]
User=polkadot
Group=polkadot
ExecStart=/usr/local/bin/polkadot \
         --name {{ project|default('project') }}-sv-public-{{ groups['public'].index(inventory_hostname) }} \
         --chain={{ chain }} \
         {% if hostvars[inventory_hostname].polkadot_additional_common_flags is defined and hostvars[inventory_hostname].polkadot_additional_common_flags|length %}
         {{ hostvars[inventory_hostname].polkadot_additional_common_flags }} \
         {% endif %}
         {% if hostvars[inventory_hostname].polkadot_additional_public_flags is defined and hostvars[inventory_hostname].polkadot_additional_public_flags|length %}
         {{ hostvars[inventory_hostname].polkadot_additional_public_flags }} \
         {% endif %}
         {% if hostvars[inventory_hostname].loggingFilter is defined and hostvars[inventory_hostname].loggingFilter|length %}
         -l{{ hostvars[inventory_hostname].loggingFilter }} \
         {% endif %}
         {% for reserved_node in groups['validator'] %}
         --sentry /ip4/{{ hostvars[reserved_node].vpnpeer_address }}/tcp/30333/p2p/{{ hostvars[reserved_node].p2p_peer_id }} \
         {% endfor %}
         {% if hostvars[inventory_hostname].telemetryUrl is defined and hostvars[inventory_hostname].telemetryUrl|length %}
         --telemetry-url '{{ hostvars[inventory_hostname].telemetryUrl }} 1'
         {% else %}
         --no-telemetry
         {% endif %}

Restart=always

[Install]
WantedBy=multi-user.target
