# ansible/roles/polkadot-public/templates/polkadot.service.j2
[Unit]
Description=Polkadot Public Node

[Service]
ExecStart=/usr/local/bin/polkadot \
         --name w3f-sv-public-{{ index }} \
         --node-key={{ node_key }} \
         --chain={{ chain }}

Restart=always

[Install]
WantedBy=multi-user.target