# ansible/roles/polkadot-common/templates/polkadot.service.j2
[Unit]
Description=Polkadot Dummy Node

[Service]
User=polkadot
Group=polkadot
ExecStart=/usr/local/bin/polkadot \
         --name polkadot-dummy

[Install]
WantedBy=multi-user.target
