# ansible/roles/polkadot-validator/templates/polkadot.service.j2
[Unit]
Description=Polkadot Node

[Service]
User=polkadot
Group=polkadot
ExecStart=/usr/local/bin/polkadot \
  --validator \
  --chain={{ chain }} \
  --listen-addr=/ip4/127.0.0.1/tcp/{{ p2p_port }} \
  --rpc-methods=Unsafe \
  {#
  ## NODE NAME
  #}
  {% if hostvars[inventory_hostname].nodeName is defined and hostvars[inventory_hostname].nodeName|length %}
  --name {{ hostvars[inventory_hostname].nodeName }} \
  {% else %}
  --name {{ project|default('project') }}-sv-validator-{{ groups['validator'].index(inventory_hostname) }} \
  {% endif %}
  {#
  ## (RUNTIME) EXECUTION
  #}
  {% if hostvars[inventory_hostname].execution is defined and hostvars[inventory_hostname].execution|length %}
  --execution {{ hostvars[inventory_hostname].execution }} \
  {% endif %}
  {#
  ## WASM EXECUTION
  #}
  {% if hostvars[inventory_hostname].wasmExecution is defined and hostvars[inventory_hostname].wasmExecution|length %}
  --wasm-execution {{ hostvars[inventory_hostname].wasmExecution }} \
  {% else %}
  --wasm-execution Compiled \
  {% endif %}
  {#
  ## TELEMETRY
  #}
  {% if hostvars[inventory_hostname].telemetryUrl is defined and hostvars[inventory_hostname].telemetryUrl|length %}
  --telemetry-url '{{ hostvars[inventory_hostname].telemetryUrl }} 1'
  {% else %}
  --no-telemetry
  {% endif %}
  {#
  ## REVERSE PROXY
  #}
  {% if hostvars[inventory_hostname].enableReverseProxy is defined and hostvars[inventory_hostname].enableReverseProxy %}
  --public-addr=/ip4/{{ hostvars[inventory_hostname].public_ip.json.ip }}/tcp/{{ proxy_port }} \
  {% endif %}
  {#
  ## LOGGING FILTER
  #}
  {% if hostvars[inventory_hostname].loggingFilter is defined and hostvars[inventory_hostname].loggingFilter|length %}
  -l{{ hostvars[inventory_hostname].loggingFilter }} \
  {% endif %}
  {#
  ## BASE/DB PATH
  #}
  {% if hostvars[inventory_hostname].basePath is defined and hostvars[inventory_hostname].basePath|length %}
  --base-path '{{ hostvars[inventory_hostname].basePath }}' \
  {% endif %}
  {#
  ## ADDITIONAL FLAGS
  #}
  {% if hostvars[inventory_hostname].additionalFlags is defined and hostvars[inventory_hostname].additionalFlags|length %}
  {{ hostvars[inventory_hostname].additionalFlags }} \
  {% endif %}
  {#
  #
  # --- DEPRECATED flags after this line. Used for backwards compatibility ---
  #
  #}
  {#
  ## WASM EXECUTION
  #}
  {% if hostvars[inventory_hostname].wasm_execution is defined and hostvars[inventory_hostname].wasm_execution|length %}
  --wasm-execution {{ hostvars[inventory_hostname].wasm_execution }} \
  {% else %}
  --wasm-execution Compiled \
  {% endif %}
  {#
  ## BASE/DB PATH
  #}
  {% if hostvars[inventory_hostname].base_path is defined and hostvars[inventory_hostname].base_path|length %}
  --base-path '{{ hostvars[inventory_hostname].base_path }}' \
  {% endif %}
  {#
  ## ADDITIONAL COMMON FLAGS
  #}
  {% if hostvars[inventory_hostname].polkadot_additional_common_flags is defined and hostvars[inventory_hostname].polkadot_additional_common_flags|length %}
  {{ hostvars[inventory_hostname].polkadot_additional_common_flags }} \
  {% endif %}
  {#
  ## ADDITIONAL FLAGS
  #}
  {% if hostvars[inventory_hostname].polkadot_additional_validator_flags is defined and hostvars[inventory_hostname].polkadot_additional_validator_flags|length %}
  {{ hostvars[inventory_hostname].polkadot_additional_validator_flags }} \
  {% endif %}

Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
