# ansible/roles/polkadot-validator/templates/polkadot.service.j2
[Unit]
Description=Polkadot Node

[Service]
User=polkadot
Group=polkadot
LimitNOFILE=10000
ExecStart=/usr/local/bin/polkadot \
  {% if hostvars[inventory_hostname].nodeName is defined and hostvars[inventory_hostname].nodeName|length %}
  --name {{ hostvars[inventory_hostname].nodeName }} \
  {% else %}
  --name {{ project|default('project') }}-sv-validator-{{ groups['validator'].index(inventory_hostname) }} \
  {% endif %}
  --validator \
  {% if hostvars[inventory_hostname].enableReverseProxy is defined and hostvars[inventory_hostname].enableReverseProxy %}
  --public-addr=/ip4/{{ hostvars[inventory_hostname].public_ip.json.ip }}/tcp/{{ proxy_port }} \
  {% endif %}
  --listen-addr=/ip4/127.0.0.1/tcp/{{ p2p_port }} \
  --rpc-methods=Unsafe \
  {% if hostvars[inventory_hostname].polkadotAdditionalCustomFlags is defined and hostvars[inventory_hostname].polkadotAdditionalCustomFlags|length %}
  {{ hostvars[inventory_hostname].polkadotAdditionalCustomFlags }} \
  {% endif %}
  {% if hostvars[inventory_hostname].polkadot_additional_common_flags is defined and hostvars[inventory_hostname].polkadot_additional_common_flags|length %}
  {{ hostvars[inventory_hostname].polkadot_additional_common_flags }} \
  {% endif %}
  {% if hostvars[inventory_hostname].polkadot_additional_validator_flags is defined and hostvars[inventory_hostname].polkadot_additional_validator_flags|length %}
  {{ hostvars[inventory_hostname].polkadot_additional_validator_flags }} \
  {% endif %}
  {% if hostvars[inventory_hostname].loggingFilter is defined and hostvars[inventory_hostname].loggingFilter|length %}
  -l{{ hostvars[inventory_hostname].loggingFilter }} \
  {% endif %}
  --chain={{ chain }} \
  {% if hostvars[inventory_hostname].base_path is defined and hostvars[inventory_hostname].base_path|length %}
  --base-path '{{ hostvars[inventory_hostname].base_path }}' \
  {% endif %}
  --wasm-execution Compiled \
  {% if hostvars[inventory_hostname].telemetryUrl is defined and hostvars[inventory_hostname].telemetryUrl|length %}
  --telemetry-url '{{ hostvars[inventory_hostname].telemetryUrl }} 1'
  {% else %}
  --no-telemetry
  {% endif %}

Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
