# ansible/roles/polkadot-validator/templates/polkadot.service.j2
{# --- PREPARE VARIABLES #}
{% set nodeName = hostvars[inventory_hostname].nodeName|default(project|default('project')+'-sv-validator-'+groups['validator'].index(inventory_hostname)) %}
{% set execution = hostvars[inventory_hostname].execution|default(None) %}
{% set wasmExecution = hostvars[inventory_hostname].wasmExecution|default('Compiled') %}
{% set telemetryUrl = hostvars[inventory_hostname].telemetryUrl|default(None) %}
{% set enableReverseProxy = hostvars[inventory_hostname].enableReverseProxy|default(false) %}
{% set loggingFilter = hostvars[inventory_hostname].loggingFilter|default(None) %}
{% set basePath = hostvars[inventory_hostname].basePath|default(None) %}
{% set additionalFlags = hostvars[inventory_hostname].additionalFlags|default(None) %}
{# ## Legacy/deprecated flags #}
{% set wasmExecution = hostvars[inventory_hostname].wasm_execution|default('Compiled') %}
{% set basePath = hostvars[inventory_hostname].base_path|default(None) %}
{% set polkadotAdditionalCommonFlags = hostvars[inventory_hostname].polkadot_additional_common_flags|default(None) %}
{% set polkadotAdditionalValidatorFlags = hostvars[inventory_hostname].polkadot_additional_validator_flags|default(None) %}
{# --- #}
[Unit]
Description=Polkadot Node

[Service]
User=polkadot
Group=polkadot
ExecStart=/usr/local/bin/polkadot \
  --validator \
  --chain={{ chain }} \
  --listen-addr=/ip4/127.0.0.1/tcp/{{ p2p_port }} \
  --rpc-methods=Unsafe \
  --name {{ nodeName }} \
  --execution {{ execution }} \
  --wasm-execution {{ wasmExecution }} \
  {% if telemetryUrl is defined %}
  --telemetry-url '{{ telemetryUrl }} 1' \
  {% else %}
  --no-telemetry \
  {% endif %}
  {% if enableReverseProxy %}
  --public-addr=/ip4/{{ hostvars[inventory_hostname].public_ip.json.ip }}/tcp/{{ proxy_port }} \
  {% endif %}
  {% if loggingFilter is defined %}
  -l{{ loggingFilter }} \
  {% endif %}
  {% if basePath is defined %}
  --base-path '{{ basePath }}' \
  {% endif %}
  {% if additionalFlags is defined %}
  {{ additionalFlags }} \
  {% endif %}
  {# --- DEPRECATED --- #}
  {% if polkadotAdditionalCommonFlags is defined %}
  {{ polkadotAdditionalCommonFlags }} \
  {% endif %}
  {% if polkadotAdditionalValidatorFlags is defined %}
  {{ polkadotAdditionalValidatorFlags }} \
  {% endif %}
  {# --- #}

Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
